# TalleySoft-Vision/firmware/rootfs-overlay/etc/systemd/system/sphere-quest-relay.service
[Unit]
Description=Duplicate Sphere Link FPV video to HDMI and Quest USB
After=sphere-quest-usb.service
Wants=sphere-quest-usb.service
Before=wifibroadcast.service wifibroadcast@gs.service

[Service]
Type=simple
ExecStart=/usr/local/bin/sphere-quest-relay
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
